# deployer/configs/supervisor.conf
# ==================================================
# SUPERVISOR CONFIGURATION FOR GESTMAT V2 QUEUES
# ==================================================
# Place this file in: /etc/supervisor/conf.d/gestmat.conf
# Update supervisor: sudo supervisorctl reread
# Add processes: sudo supervisorctl update
# Start: sudo supervisorctl start gestmat-queue:*
# ==================================================

# ==================================================
# QUEUE WORKER - Default Queue
# ==================================================
[program:gestmat-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/gestmatv2/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --timeout=300
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gestmat-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

# ==================================================
# QUEUE WORKER - High Priority (PDF Generation)
# ==================================================
[program:gestmat-queue-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/gestmatv2/artisan queue:work redis --queue=high,default --sleep=1 --tries=3 --max-time=1800 --timeout=300
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gestmat-queue-high.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

# ==================================================
# QUEUE WORKER - Low Priority (Imports/Exports)
# ==================================================
[program:gestmat-queue-low]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/gestmatv2/artisan queue:work redis --queue=low --sleep=5 --tries=3 --max-time=7200 --timeout=600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gestmat-queue-low.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=7200

# ==================================================
# SCHEDULER (Laravel Schedule Runner)
# ==================================================
[program:gestmat-scheduler]
process_name=%(program_name)s
command=/bin/sh -c "while [ true ]; do (php /var/www/gestmatv2/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gestmat-scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3

# ==================================================
# GROUP CONFIGURATION
# ==================================================
[group:gestmat]
programs=gestmat-queue,gestmat-queue-high,gestmat-queue-low,gestmat-scheduler
priority=999

# ==================================================
# USAGE COMMANDS
# ==================================================
# Status:      sudo supervisorctl status gestmat:*
# Start all:   sudo supervisorctl start gestmat:*
# Stop all:    sudo supervisorctl stop gestmat:*
# Restart all: sudo supervisorctl restart gestmat:*
# Logs:        sudo supervisorctl tail -f gestmat-queue
# Update:      sudo supervisorctl reread && sudo supervisorctl update
#
# Individual control:
# sudo supervisorctl start gestmat-queue:gestmat-queue_00
# sudo supervisorctl stop gestmat-queue-high:*
# sudo supervisorctl restart gestmat-scheduler
#
# ==================================================
# NOTES
# ==================================================
# 1. Adjust path /var/www/gestmatv2 to your installation path
# 2. Queue connection should match .env (redis)
# 3. numprocs determines how many worker processes run
# 4. --max-time prevents memory leaks (workers restart)
# 5. stopwaitsecs should be >= max-time to allow graceful shutdown
# 6. Monitor logs: tail -f /var/log/supervisor/gestmat-*.log
# 7. After deployment: sudo supervisorctl restart gestmat:*
